# {{cookiecutter.out_dir}}/sys/etc/cron.d/drupal_cron.frep
MAILTO=""
# Launching a cron scheduler
* * * * *   drupal  [ ! -e /code/app/var/private/SUSPEND_CRONS ] && /code/app/bin/drush core:cron
